# systemd/email-reader.timer
[Unit]
Description=Gmail-to-PDF Email Reader Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=__POLL_INTERVAL__min

[Install]
WantedBy=timers.target
